# playbooks/roles/oc-mentor-assistant/templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket for {{ app_name }} app

[Socket]
ListenStream=/run/{{ app_name }}.sock

[Install]
WantedBy=sockets.target